# fancontrol.service
[Unit]
Description=Fan Control Script

[Service]
ExecStart=/usr/bin/fancontrol.sh
Restart=always

[Install]
WantedBy=multi-user.target
